# cutedog.conf
[program:dogweb]
command=/usr/local/bin/gunicorn wsgi --bind 0.0.0.0:8002 --pid /tmp/cutedog.pid
command=export ENV_KOHI="pro"
directory=/var/www/cutedog
autostart=true
autorestart=true